# Library/LaunchAgents/us.theodells.kopia.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>us.theodells.restic.backup</string>

    <key>WorkingDirectory</key>
    <string>/Users/michael.odell</string>

    <key>StandardErrorPath</key>
    <string>/Users/michael.odell/log/kopia-backup.log</string>

    <key>StandardOutPath</key>
    <string>/Users/michael.odell/log/kopia-backup.log</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/michael.odell/bin/kopia-backup</string>
    </array>

    <key>StartInterval</key>
    <integer>1800</integer> <!-- Every 30 minutes -->

    <key>RunAtLoad</key>
    <true/>

    <key>ProcessType</key>
    <string>Background</string>

    <key>LowPriorityIO</key>
    <true/>

    <key>LowPriorityBackgroundIO</key>
    <true/>

    <key>MaterializeDatalessFiles</key>
    <true/>

</dict>
</plist>
